# rootfs/usr/share/hoth-os/apps/arr-stack/quadlets/flaresolverr.container
[Unit]
Description=FlareSolverr - Proxy to bypass Cloudflare protections
After=network-online.target

[Container]
Image=ghcr.io/flaresolverr/flaresolverr:latest
AutoUpdate=registry
Pod=arr-stack.pod
Environment=LOG_LEVEL=info

HealthCmd=/usr/bin/curl --fail http://localhost:8191/health
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3

HealthStartupCmd=/usr/bin/curl --fail --retry 3 --max-time 5 http://localhost:8191/health
HealthStartupInterval=5s
HealthStartupTimeout=5s
HealthStartupRetries=30

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target