# rootfs/etc/profile.d/hoth-os.sh
alias hjust="just --justfile /usr/share/hoth-os/justfile"
